# systemd/user/gopls@.service
[Unit]
Description=Go Language Server Per-Connection Server

[Service]
ExecStart=-/home/pablo/go/bin/gopls serve
StandardInput=socket
